# on.timer
[Unit]
Description=Start timer

[Timer]
OnCalendar=Mon..Fri *-*-* 18:20:00
Unit=on.target
Persistent=true

[Install]
WantedBy=timers.target
